# debianpkgs/audit/debian/auditd.postinst
#!/bin/sh

set -e

case "$1" in
    configure)
	# Remove old rc links.
  	if dpkg --compare-versions "$2" lt "1.6.5-0ubuntu2"; then
    		update-rc.d -f auditd remove > /dev/null
  	fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

pam-auth-update --package

#DEBHELPER#

exit 0
