# packaging/postfix/debian/postfix.service
[Unit]
Description=Postfix mail transfer agent
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/postfix -c /etc/postfix start
ExecStop=/usr/sbin/postfix -c /etc/postfix stop
ExecReload=/usr/sbin/postfix -c /etc/postfix reload

[Install]
WantedBy=multi-user.target
