# packaging/isc-dhcp/debian/isc-dhcp-relay.prerm
#!/bin/sh -e
#
#

invoke-rc.d isc-dhcp-relay stop

#DEBHELPER#

exit 0
