# packaging/apt-show-versions/debian/cron.daily
#!/bin/sh

[ -f /usr/bin/apt-show-versions ] || exit 0

apt-show-versions -i
