# packaging/slrn/debian/slrn.cron.weekly
#!/bin/sh
#
# Cron.weekly script for slrn. Update newsgroups descriptions.

[ -x /usr/bin/slrn ] || exit 0

if [ -f /etc/default/slrn ]; then
	. /etc/default/slrn
fi

if [ "$GETDESCS" = "cron job" ]; then
	/usr/sbin/slrn_getdescs $INTERVAL >/dev/null 2>&1
fi
