# packaging/nas/debian/nas.prerm
#!/bin/sh -e

if [ \( "$1" = "upgrade" -o "$1" = "remove" \) -a -L /usr/doc/nas ]; then
 	rm -f /usr/doc/nas
fi

case "$1" in
	remove|upgrade|remove-in-favour|deconfigure-in-favour)


		# Stop running nas
        if which invoke-rc.d >/dev/null 2>&1; then
            invoke-rc.d nas stop
        else
            /etc/init.d/nas stop
        fi
		;;

	failed-upgrade)
		;;
esac
