# packaging/nas/debian/nas.conffiles
/etc/nas/nasd.conf
/etc/init.d/nas
/etc/default/nas
