# packaging/xen/debian/xen-utils-common.preinst
#!/bin/sh

set -e

case "$1" in
install|upgrade)
	;;

abort-upgrade)
	;;

*)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
	;;
esac

dpkg-maintscript-helper mv_conffile /etc/init.d/xend /etc/init.d/xen 4.1.2-4~ -- "$@"
update-rc.d -f xend remove >/dev/null

#DEBHELPER#

exit 0
