# packaging/nagios3/debian/nagios3-core.postinst
#!/bin/sh -e

#if we stop nagios3 in nagios3.prerm we should also start it here..
#(#481334)
if [ -x "/etc/init.d/nagios3" ]; then
	if ! /etc/init.d/nagios3 status >/dev/null; then

		if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
			invoke-rc.d nagios3 start || true
		else
			/etc/init.d/nagios3 start || true
		fi
	fi
fi

#DEBHELPER#
