# packaging/lm-sensors/debian/lm-sensors.service
[Unit]
Description=Initialize hardware monitoring sensors
After=syslog.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/sensors -s
ExecStart=/usr/bin/sensors

[Install]
WantedBy=multi-user.target
