# packaging/krb/debian/krb5-kdc-ldap.insserv-override
### BEGIN INIT INFO
# Provides:             krb5-kdc
# Required-Start:       $local_fs $remote_fs $network $syslog
# Required-Stop:        $local_fs $remote_fs $network $syslog
# Should-Start:         slapd
# Should-Stop:          slapd
# X-Start-Before: $x-display-manager                                          
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
### END INIT INFO
