# packaging/redhat-cluster/debian/cman.postinst
#!/bin/sh

set -e

# Those using dependency based boot sequencing with sysv-rc and
# installing cman before version 2.20081102-1.1 would have wrong
# runlevel symlinks.  Recover from this.
if [ "$1" = "configure" ] && dpkg --compare-versions "$2" le "3.0.2-5" \
   && [ -f /etc/rc2.d/S[0-9][0-9]cman ] && [ ! -f /etc/rcS.d/S[0-9][0-9]cman ]
then
    update-rc.d -f cman remove
fi

#DEBHELPER#
