# packaging/ifupdown/debian/ifupdown-clean.service
[Unit]
Description=Clean old interface status during boot
DefaultDependencies=no
After=network.service
After=remount-rootfs.service

[Service]
Type=oneshot
ExecStart=-/bin/rm -f /etc/network/run/ifstate
RemainAfterExit=true
