# packaging/ifupdown/debian/ifup@.service
[Unit]
Description=Start ifup for %I
After=local-fs.target
After=network.service
After=ifupdown-clean.service
After=ifupdown-mkdir.service
Wants=ifupdown-clean.service
Wants=ifupdown-mkdir.service

[Service]
ExecStart=/sbin/ifup --allow=hotplug %I
ExecStop=/sbin/ifdown %I
RemainAfterExit=true
