# packaging/uwsgi/debian/libapache2-mod.postinst.in
#!/bin/sh
# postinst script for libapache2-mod-@@module_name@@
#
# see: dh_installdeb(1)

set -e

apache_force_reload() {
  if apache2ctl configtest 2>/dev/null; then
    invoke-rc.d apache2 force-reload || true
  else
    echo "Your apache2 configuration is broken, please fix it and restart apache2 manually."
  fi
}

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
      if [ -z "$2" ]; then # installing package
        if [ -e /etc/apache2/apache2.conf ]; then
          a2enmod @@module_name@@ >/dev/null || true
          apache_force_reload
        fi
      else # upgrading package
        if [ -e /etc/apache2/mods-enabled/@@module_name@@.load ]; then
          apache_force_reload
        fi
      fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
