# packaging/openssh/debian/openssh-server.service
[Unit]
Description=OpenBSD SSH server
After=network.target

[Service]
EnvironmentFile=/etc/default/ssh
ExecStart=/usr/sbin/sshd $SSHD_OPTS -D
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
