# packaging/openntpd/debian/openntpd.service
[Unit]
Description=OpenBSD NTP server
Requires=network.target

[Service]
EnvironmentFile=/etc/default/openntpd
ExecStart=/usr/sbin/ntpd $DAEMON_OPTS -d
ExecStop=killall ntpd
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
