# packaging/openntpd/debian/openntpd.service
[Unit]
Description=ntpd
Requires=network.target

[Service]
ExecStart=/usr/sbin/ntpd -s

[Install]
WantedBy=multi-user.target
