# packaging/virtualbox/debian/virtualbox.prerm
#!/bin/sh

set -e

if [ -x "/etc/init.d/virtualbox" ] && [ "$1" = remove ]; then
	invoke-rc.d virtualbox stop || true
fi

#DEBHELPER#
