# packaging/libpam-ssh/debian/postinst
#!/bin/sh
set -e

case "$1" in

    configure)
        if [ -f /etc/init.d/libpam-ssh ]; then
            echo -n "Removing trailing old script.."
            update-rc.d -f libpam-ssh remove > /dev/null 2>&1
            echo -n "."
            rm -f /etc/init.d/libpam-ssh > /dev/null 2>&1
            echo "Ok."
        fi
        ## Install setup files in /usr/share/pam-configs.
        pam-auth-update --package
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ## Do nothing.
	;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
	;;
esac

#DEBHELPER#
