# packaging/packagekit/debian/packagekit.postrm
#!/bin/sh
# postrm script for packagekit
set -e

if [ "$1" = "purge" ]; then

        if [ -x /etc/init.d/dbus ]; then
                invoke-rc.d dbus force-reload || true
        fi

        # Remove PK transaction db and desktop cache
        rm -rf /var/lib/PackageKit
        rm -rf /var/cache/PackageKit

fi

#DEBHELPER#
