# packaging/desktop-base/debian/postinst
#!/bin/bash -e
set -e

#DEBHELPER#

if dpkg-maintscript-helper supports rm_conffile 2>/dev/null; then
  dpkg-maintscript-helper rm_conffile /etc/kde3/kdeglobals 6.0.1 -- "$@"
fi

if [ "$1" = configure ] && dpkg --compare-versions "$2" lt-nl 7; then
  update-alternatives --remove desktop-grub /usr/share/images/desktop-base/moreblue-orbit-grub.png
fi

if [ "$1" = "configure" -o "$1" = "abort-upgrade" ]; then
    # Alternatives for the background
    while read background priority; do
        update-alternatives --install \
            /usr/share/images/desktop-base/desktop-background \
            desktop-background \
            /usr/share/images/desktop-base/$background $priority
    done << EOF
joy-wallpaper_1600x1200.svg 65
joy-wallpaper_1280x1024.svg 65
joy-wallpaper_1920x1080.svg 70
joy-wallpaper_1920x1200.svg 65
joy-inksplat-wallpaper_1920x1080.svg 65
spacefun-wallpaper.svg 60
spacefun-wallpaper-widescreen.svg  55
moreblue-orbit-wallpaper.svg 50
moreblue-orbit-wallpaper-widescreen.svg 50
nightly-wallpaper.png 40
debian-blueish-wallpaper.svg 30
EOF

    # Set up an alternative for the XML version of the background
    # (for GNOME)
    update-alternatives --install \
        /usr/share/images/desktop-base/desktop-background.xml \
        desktop-background.xml \
        /usr/share/images/desktop-base/joy.xml 10

    # Alternatives for the splash
    while read splash priority; do
        update-alternatives --install \
            /usr/share/images/desktop-base/desktop-splash \
            desktop-splash \
            /usr/share/images/desktop-base/$splash $priority
    done << EOF
spacefun-splash.svg 50
moreblue-orbit-splash.png 40
gnome-splash-curves.png 30
EOF

    # Alternatives for grub
    while read background priority; do
        update-alternatives --install \
            /usr/share/images/desktop-base/desktop-grub.png \
            desktop-grub \
            /usr/share/images/desktop-base/$background $priority
    done << EOF
joy-grub.png 20
spacefun-grub.png 15
spacefun-grub-widescreen.png 14
EOF

    # GRUB background
    if which update-grub2 > /dev/null ; then
        sync
        update-grub2 || true
    fi

    if [ -x /usr/sbin/update-initramfs ]; then
        update-initramfs -u
    fi

    # Reload gdm3 configuration
    if [ -x /usr/sbin/gdm3 ]; then
        invoke-rc.d gdm3 reload || true
    fi

fi
