# packaging/lldpad/debian/lldpad.postinst
#!/bin/sh

set -e

if [ -x "/etc/init.d/lldpad" ]; then
	update-rc.d lldpad start 20 2 3 4 5 . stop 80 0 1 6 . >/dev/null
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d lldpad start || exit $?
	else
		/etc/init.d/lldpad start || exit $?
	fi
fi

#DEBHELPER#

exit 0
