# packaging/lldpad/debian/lldpad.postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ] ; then
	update-rc.d lldpad remove >/dev/null
fi

#DEBHELPER#

exit 0
