# packaging/openslp/debian/postinst
#!/bin/sh

set -ex

# update init script symlinks with new priorities when upgrading
# from older versions
if [ "$1" = configure ] && dpkg --compare-versions "$2" le 1.2.1-7.6; then
	echo "Reinstalling init script for new priorities ..." >&2
	update-rc.d slpd remove >/dev/null
fi

#DEBHELPER#
