# packaging/gpm/debian/gpm.service
[Unit]
Description=General Protocol Mouse daemon
DefaultDependencies=no
ConditionPathExists=/dev/tty0
ConditionPathExists=/dev/input/mice

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/gpm.conf
ExecStart=/usr/sbin/gpm -m $device -t $type

[Install]
WantedBy=getty.target
