# packaging/systemd/debian/systemd-sysv.install
# udevadm.8 is installed by systemd (main package)
usr/share/man/man1/init.1
usr/share/man/man8/runlevel.8
usr/share/man/man8/shutdown.8
usr/share/man/man8/poweroff.8
usr/share/man/man8/reboot.8
usr/share/man/man8/halt.8
usr/share/man/man8/telinit.8
usr/share/man/man8/systemd-*.8
debian/invoke-rc.d usr/sbin/
debian/update-rc.d usr/sbin/
