# packaging/zookeeper/debian/zookeeper.cron.daily
#!/bin/sh

ZK_PURGE_SCRIPT=/usr/share/zookeeper/bin/zkPurgeTxnLog.sh
[ -r $ZK_PURGE_SCRIPT ] && . $ZK_PURGE_SCRIPT || :
