# packaging/git/debian/git-daemon-sysvinit.prerm
#!/bin/sh
set -e

test "$1" = 'remove' || test "$1" = 'deconfigure' || \
  test "$1" = 'failed-upgrade' || exit 0

if [ -x "/etc/init.d/git-daemon" ]; then
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d git-daemon stop || exit $?
	else
		/etc/init.d/git-daemon stop || exit $?
	fi
fi
