# packaging/git/debian/git-daemon-sysvinit.conffiles
/etc/default/git-daemon
/etc/init.d/git-daemon
