# debianpkgs/libvirt/debian/libvirt-sanlock.cron.weekly
#!/bin/sh

# Disabled by default
exit 0

/usr/sbin/virt-sanlock-cleanup -q
exit 0



