# debianpkgs/virtualbox/debian/virtualbox-dkms.postinst
#!/bin/sh

set -e

#DEBHELPER#

if [ "$1" = configure ]; then
	# only restart if VirtualBox isn't running
	if test -x /etc/init.d/virtualbox && ! pidof VBoxSVC > /dev/null; then
		invoke-rc.d virtualbox restart || true
	fi
fi
