# packaging/console-setup/setupcon.service
[Unit]
Description=Set up the console
DefaultDependencies=no
ConditionPathExists=/dev/tty0

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/setupcon
EnvironmentFile=/etc/default/console-setup

[Install]
WantedBy=getty.target
