# packaging/console-setup/setupcon.service
[Unit]
Description=Set up the console
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target
ConditionPathExists=/dev/tty0

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/setupcon
