# packaging/alsa-utils/debian/udev.rules
KERNEL=="controlC[0-7]", ACTION=="add", RUN+="/lib/udev/alsa-utils"

# Handle addition and removal of soundcards.
KERNEL=="controlC*", ACTION=="add", SUBSYSTEM=="sound", \
  RUN+="/usr/sbin/alsactl restore /dev/$name"
KERNEL=="controlC*", ACTION=="remove", SUBSYSTEM=="sound", \
  RUN+="/usr/sbin/alsactl store /dev/$name"
